# chadrouter/etc/systemd/system/update_ddns.timer
[Unit]
Description=Run Update Dynamic DNS
After=network-online.target

[Timer]
OnBootSec=30s
OnUnitInactiveSec=1h

[Install]
WantedBy=timers.target
