# any/etc/systemd/system/sshd.service.d/override.conf
[Unit]
Wants=sshdgenkeys.service network-online.target
After=sshdgenkeys.service network.target network-online.target

[Service]
Restart=always
RestartSec=5s